# crontab
*  *  *  *  * /usr/bin/php /path/to/your/owncloud/occ system:cron
